# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash
#
# Cron
#

if [ -z "$USE_CRON" ]; then
  exit 0
fi

crond -f
