# root/etc/periodic/daily/10-crlupdate.sh
#!/usr/bin/with-contenv bash
#
# CRL update
#
ovpn pki crl
