# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

CONFIGS=$(ls -1 /config/openvpn/server | wc -l) # Get number of files in /config/openvpn/server
if [ $CONFIGS -gt 0 ]; then
    # Run OpenVPN server only if configs are present
    echo "Running OpenVPN ..."
    sudo -E -u abc openvpn --cd /config --config ${OVPN_RUN}
else
    # No OpenVPN config present
    echo "No OpenVPN config present"
    
    # Do nothing until container restart
    while true
    do
        sleep 3600
    done
fi