# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

openvpn --cd /config --config ${OVPN_RUN}