# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

sudo -E -u abc openvpn --cd /config --config ${OVPN_RUN}